# mod_for_dbusspy/dbus.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# Drop this in as dbus.conf after renaming the normal dbus.conf to
# "realdbus.conf". That will ensure this attaches to dbus first.

start on started realdbus

respawn

exec dbus-spy -w /var/log/dbusspy.log
